# src/linux/bing-wallpaper.service
[Unit]
Description=Bing daily wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/bing-wallpaper.sh -w -L -m pt-BR -r UHD
